# test-env/rootfs/etc/init.d/firewall
#!/sbin/openrc-run

name="firewall"

depend() {
    after net
}

start() {
    ebegin "Starting firewall"
    /sbin/nft -f /etc/nftables.conf 2>/dev/null
    eend $?
}

stop() {
    ebegin "Stopping firewall"
    /sbin/nft flush ruleset 2>/dev/null
    eend $?
}
