# test-env/rootfs/etc/init.d/webui
#!/sbin/openrc-run
# Web UI (lighttpd) service for AX6600

name="webui"
description="Web management interface"

LIGHTTPD_CONF="/etc/lighttpd/lighttpd.conf"

depend() {
    after net
    need net
}

start() {
    ebegin "Starting Web UI"
    
    # Check if lighttpd is installed
    if ! command -v lighttpd > /dev/null; then
        ewend 1 "lighttpd not installed (apk add lighttpd)"
        return 1
    fi
    
    # Create config if not exists
    if [ ! -f "$LIGHTTPD_CONF" ]; then
        create_lighttpd_config
    fi
    
    # Start lighttpd
    lighttpd -f "$LIGHTTPD_CONF" -D &
    
    eend $?
}

stop() {
    ebegin "Stopping Web UI"
    
    killall lighttpd 2>/dev/null || true
    
    eend $?
}

create_lighttpd_config() {
    mkdir -p /etc/lighttpd
    
    cat > "$LIGHTTPD_CONF" << 'EOF'
# Lighttpd configuration for AX6600 Web UI

server.document-root = "/www"
server.port = 80
server.username = "root"
server.groupname = "root"

# Logging
server.errorlog = "/var/log/lighttpd/error.log"
accesslog.filename = "/var/log/lighttpd/access.log"

# Index files
index-file.names = ( "index.html", "index.htm" )

# CGI
server.modules = (
    "mod_cgi",
    "mod_accesslog",
    "mod_auth",
    "mod_authn_file"
)

# CGI handler
cgi.assign = (
    ".cgi" => "/bin/sh",
    ".sh" => "/bin/sh"
)

# Directory listings
dir-listing.activate = "disable"

# Authentication (optional)
# auth.backend = "htpasswd"
# auth.backend.htpasswd.userfile = "/etc/lighttpd/.htpasswd"
# auth.require = ( "" => (
#     "method" => "basic",
#     "realm" => "AX6600",
#     "require" => "valid-user"
# ))

# MIME types
mimetype.assign = (
    ".html" => "text/html",
    ".htm" => "text/html",
    ".css" => "text/css",
    ".js" => "application/javascript",
    ".json" => "application/json",
    ".png" => "image/png",
    ".jpg" => "image/jpeg",
    ".gif" => "image/gif",
    ".ico" => "image/x-icon"
)
EOF
    
    mkdir -p /var/log/lighttpd
    
    echo "Lighttpd configuration created"
}

status() {
    echo "=== Web UI Status ==="
    
    if command -v lighttpd > /dev/null; then
        echo "lighttpd: installed"
        
        if pgrep lighttpd > /dev/null; then
            echo "Service: running on port 80"
            echo "URL: http://$(ip addr show br-lan | grep 'inet ' | awk '{print $2}' | cut -d/ -f1)/"
        else
            echo "Service: not running"
        fi
    else
        echo "lighttpd: not installed"
        echo "Install: apk add lighttpd"
    fi
}