# test-env/rootfs/etc/init.d/wifi
#!/sbin/openrc-run

name="WiFi"

depend() {
    need net
    after bootmisc
}

start() {
    ebegin "Starting WiFi"
    
    # Load WiFi firmware
    modprobe ath11k 2>/dev/null || true
    
    # Wait for interfaces
    sleep 2
    
    # Start hostapd for 2.4GHz
    if [ -f /etc/hostapd/hostapd.conf ]; then
        /usr/sbin/hostapd -B /etc/hostapd/hostapd.conf 2>/dev/null || true
    fi
    
    # Start hostapd for 5GHz
    if [ -f /etc/hostapd/hostapd-5g.conf ]; then
        /usr/sbin/hostapd -B /etc/hostapd/hostapd-5g.conf 2>/dev/null || true
    fi
    
    eend $?
}

stop() {
    ebegin "Stopping WiFi"
    killall hostapd 2>/dev/null || true
    eend $?
}
