# clawui/apps/kms/kms-server.init
#!/sbin/openrc-run
# KMS Server init script

name="KMS Server"
command="/usr/share/kms/kms-server.sh"
command_background="true"
pidfile="/var/run/kms-server.pid"
command_output="/var/log/kms-server.log"

depend() {
    need net
    use logger
}

start_pre() {
    checkpath -d /var/log -m 0755
    checkpath -f /var/log/kms-server.log -m 0644
}
