# test-env/rootfs/usr/sbin/mac-filter
#!/sbin/openrc-run
# MAC address filtering service for JDCloud AX6600

name="mac-filter"
description="MAC address whitelist/blacklist"

MAC_FILTER_CONF="/etc/mac-filter.conf"

depend() {
    after firewall
    need net
}

start() {
    ebegin "Starting MAC filter"
    
    if [ ! -f "$MAC_FILTER_CONF" ]; then
        touch "$MAC_FILTER_CONF"
    fi
    
    apply_mac_filter
    
    eend $?
}

stop() {
    ebegin "Stopping MAC filter"
    
    # Remove filter rules
    nft delete table inet mac-filter 2>/dev/null || true
    
    eend $?
}

apply_mac_filter() {
    local mode=$(grep "^MODE=" "$MAC_FILTER_CONF" | cut -d= -f2)
    mode="${mode:-whitelist}"  # Default to whitelist
    
    # Create nftables table for MAC filtering
    nft add table inet mac-filter 2>/dev/null || true
    nft add chain inet mac-filter input { type filter hook input priority 1 \; } 2>/dev/null || true
    nft add chain inet mac-filter forward { type filter hook forward priority 1 \; } 2>/dev/null || true
    
    # Get MAC addresses
    local macs=$(grep "^MAC=" "$MAC_FILTER_CONF" | cut -d= -f2)
    
    if [ "$mode" = "whitelist" ]; then
        # Whitelist mode: only allow listed MACs
        echo "Applying whitelist mode..."
        
        for mac in $macs; do
            nft add rule inet mac-filter forward ether saddr $mac accept 2>/dev/null || true
        done
        
        # Block all others on LAN
        nft add rule inet mac-filter forward iifname br-lan drop 2>/dev/null || true
        
    elif [ "$mode" = "blacklist" ]; then
        # Blacklist mode: block listed MACs
        echo "Applying blacklist mode..."
        
        for mac in $macs; do
            nft add rule inet mac-filter forward ether saddr $mac drop 2>/dev/null || true
        done
    fi
    
    echo "MAC filter applied (${mode} mode)"
}

status() {
    echo "=== MAC Filter Status ==="
    
    if [ -f "$MAC_FILTER_CONF" ]; then
        echo "Configuration: $MAC_FILTER_CONF"
        echo "Mode: $(grep MODE $MAC_FILTER_CONF | cut -d= -f2 || echo 'whitelist')"
        echo "MACs:"
        grep "^MAC=" "$MAC_FILTER_CONF" | cut -d= -f2 | while read mac; do
            echo "  $mac"
        done
    else
        echo "Configuration: not found"
    fi
    
    echo ""
    echo "Active rules:"
    nft list table inet mac-filter 2>/dev/null | grep ether || echo "  No rules active"
}