# rootfs-overlay/etc/init.d/network
#!/sbin/openrc-run

depend() {
    need net
    after bootmisc
}

start() {
    ebegin "Starting network"
    ifup -a 2>/dev/null || true
    eend $?
}

stop() {
    ebegin "Stopping network"
    ifdown -a 2>/dev/null || true
    eend $?
}
