# rootfs-overlay/etc/init.d/upnp
#!/sbin/openrc-run
# UPnP service for JDCloud AX6600

name="upnp"
description="UPnP port mapping service"

MINIUPNPD_CONF="/etc/miniupnpd.conf"

depend() {
    after net firewall
    need net
}

start() {
    ebegin "Starting UPnP service"
    
    # Check if miniupnpd is installed
    if ! command -v miniupnpd > /dev/null; then
        ewend 1 "miniupnpd not installed (apk add miniupnpd)"
        return 1
    fi
    
    # Create config if not exists
    if [ ! -f "$MINIUPNPD_CONF" ]; then
        create_upnp_config
    fi
    
    # Start miniupnpd
    miniupnpd -f "$MINIUPNPD_CONF" 2>/dev/null
    
    eend $?
}

stop() {
    ebegin "Stopping UPnP service"
    
    killall miniupnpd 2>/dev/null || true
    
    eend $?
}

create_upnp_config() {
    cat > "$MINIUPNPD_CONF" << 'EOF'
# miniupnpd configuration for JDCloud AX6600

# WAN interface (external)
ext_ifname=eth0

# LAN interface (internal)
listening_ip=192.168.1.1
internal_ifname=br-lan

# Port range for UPnP mapping
min_port=1024
max_port=65535

# NAT rules chain name (iptables)
# chain=MINIUPNPD

# Enable NAT-PMP (Apple)
enable_natpmp=yes

# Enable UPnP
enable_upnp=yes

# UPnP IGD version
upnp_igd_version=2

# SSDP packet lifetime (seconds)
ssdp_interval=60

# UPnP notify interval (seconds)
notify_interval=30

# Allow rules (permission)
# Allow specific IPs or deny
# allow 192.168.1.100-192.168.1.200

# Deny rules
# deny 192.168.1.50

# Secure mode (only allow specific clients)
secure_mode=no

# System uptime script
uptime_script=/usr/sbin/upnp-uptime

# Log level
log_level=3
EOF
    
    echo "UPnP configuration created"
}

status() {
    echo "=== UPnP Status ==="
    
    if command -v miniupnpd > /dev/null; then
        echo "miniupnpd: installed"
        
        if pgrep miniupnpd > /dev/null; then
            echo "Service: running"
            
            # Show port mappings if available
            echo "Port mappings:"
            miniupnpd -g 2>/dev/null | head -20 || echo "  Unable to get mappings"
        else
            echo "Service: not running"
        fi
    else
        echo "miniupnpd: not installed"
        echo "Install: apk add miniupnpd"
    fi
}