# rootfs-overlay/etc/init.d/usb
#!/sbin/openrc-run
# USB storage service for JDCloud AX6600

name="usb"
description="USB storage support"

depend() {
    after bootmisc
}

start() {
    ebegin "Starting USB storage support"
    
    # Load USB modules
    modprobe usb-storage 2>/dev/null || true
    modprobe usbhid 2>/dev/null || true
    
    # Wait for USB devices
    sleep 2
    
    # Auto-mount USB storage
    automount_usb
    
    eend $?
}

stop() {
    ebegin "Stopping USB storage"
    
    # Unmount all USB drives
    for mount_point in /mnt/usb/*; do
        if [ -d "$mount_point" ]; then
            umount "$mount_point" 2>/dev/null || true
        fi
    done
    
    eend $?
}

automount_usb() {
    mkdir -p /mnt/usb
    
    # Find all USB block devices
    for dev in /dev/sd[a-z] /dev/sd[a-z][0-9]*; do
        if [ -b "$dev" ]; then
            local name=$(basename "$dev")
            local mount_point="/mnt/usb/${name}"
            
            # Skip if already mounted
            if ! mountpoint -q "$mount_point" 2>/dev/null; then
                mkdir -p "$mount_point"
                
                # Try to mount with various filesystems
                mount -t auto "$dev" "$mount_point" 2>/dev/null || \
                mount -t vfat "$dev" "$mount_point" 2>/dev/null || \
                mount -t ntfs-3g "$dev" "$mount_point" 2>/dev/null || \
                mount -t ext4 "$dev" "$mount_point" 2>/dev/null || \
                rmdir "$mount_point" 2>/dev/null
            fi
        fi
    done
}

status() {
    echo "=== USB Storage Status ==="
    
    # Show USB devices
    echo "USB devices:"
    ls -la /dev/sd* 2>/dev/null || echo "  No USB storage devices"
    
    echo ""
    echo "Mounted:"
    mount | grep /mnt/usb || echo "  None"
    
    echo ""
    echo "USB modules:"
    lsmod | grep usb || echo "  None loaded"
}