# rootfs-overlay/etc/inittab
# /etc/inittab - Alpine Linux init configuration

::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

# Console on serial port (for debugging)
::respawn:/sbin/getty -L ttyMSM0 115200 vt100

# Local console
# tty1::respawn:/sbin/getty 38400 tty1

# Ctrl+Alt+Del action
::ctrlaltdel:/sbin/reboot

# Shutdown
::shutdown:/sbin/openrc shutdown
